# packaging/macos/installer/com.heimdex.agent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.heimdex.agent</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/Applications/Heimdex Agent.app/Contents/MacOS/heimdex-agent</string>
    </array>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>
    
    <key>StandardOutPath</key>
    <string>/tmp/heimdex-agent.stdout.log</string>
    
    <key>StandardErrorPath</key>
    <string>/tmp/heimdex-agent.stderr.log</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>HEIMDEX_LOG_LEVEL</key>
        <string>info</string>
    </dict>
    
    <key>ProcessType</key>
    <string>Interactive</string>
    
    <key>LimitLoadToSessionType</key>
    <string>Aqua</string>
</dict>
</plist>
